# live.service
[Unit]
Description=Raspberry Pi Live Stream Youtube
[Service]
ExecStart=/opt/live.sh
[Install]
WantedBy=default.target
